# supervisord.conf
; supervisord configuration for the Synology MCP Hub.
;
; Each MCP runs as its own child process on a distinct internal port.
; Logs are streamed to stdout/stderr so `docker logs mcp-hub` shows
; everything (per-program log lines are prefixed with the program name).
;
; Programs:
;   arr         - python stdio MCP wrapped by supergateway (since mcp-arr
;                 hasn't been ported to native HTTP transport)
;   dispatcharr - python FastMCP, native HTTP transport
;   gramps      - node TS MCP, native HTTP transport

[supervisord]
nodaemon=true
loglevel=%(ENV_SUPERVISOR_LOGLEVEL)s
logfile=/dev/stdout
logfile_maxbytes=0
pidfile=/run/supervisord.pid

[unix_http_server]
file=/run/supervisor.sock

[supervisorctl]
serverurl=unix:///run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

; -------------------------------------------------------------------------
; arr-mcp (Sonarr / Radarr / Prowlarr wrapper)
;
; mcp-arr installs a console script `arr-mcp` that runs the stdio MCP.
; supergateway bridges stdio ↔ streamable HTTP and exposes /mcp on 8120.
; Env (SONARR_URL, SONARR_API_KEY, RADARR_*, PROWLARR_*) is supplied by
; docker-compose env_file: env/arr.env.
; -------------------------------------------------------------------------
[program:arr]
command=supergateway
    --stdio "arr-mcp"
    --outputTransport streamableHttp
    --streamableHttpPath /mcp
    --healthEndpoint /healthz
    --port 8120
    --logLevel info
    --cors
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=10
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
stdout_events_enabled=true
stderr_events_enabled=true
environment=PORT="8120"

; -------------------------------------------------------------------------
; dispatcharr-mcp (TV channel manager wrapper)
;
; Native HTTP via FastMCP — no supergateway needed. server.py reads PORT
; and serves /mcp directly. Env (DISPATCHARR_URL, DISPATCHARR_API_KEY,
; STREAM_PROBE_*) supplied by docker-compose.
; -------------------------------------------------------------------------
[program:dispatcharr]
command=python /mcps/dispatcharr/server.py
directory=/mcps/dispatcharr
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=10
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
environment=PORT="8121",MCP_ENDPOINT="/mcp"

; -------------------------------------------------------------------------
; gramps-mcp (Gramps Web genealogy app wrapper)
;
; Native HTTP via Node TS app. dist/index.js reads PORT and serves /mcp.
; Env (GRAMPSWEB_URL, GRAMPS_USERNAME, GRAMPS_PASSWORD) supplied by
; docker-compose.
; -------------------------------------------------------------------------
[program:gramps]
command=node /mcps/gramps/dist/index.js
directory=/mcps/gramps
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=10
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
environment=PORT="8122",MCP_TRANSPORT="http"

; -------------------------------------------------------------------------
; trakt-mcp (Trakt.tv API wrapper)
;
; stdio MCP wrapped by supergateway. Trakt has its own per-user device-code
; OAuth (separate from the Worker's OAuth) — token persists to
; TRAKT_AUTH_TOKEN_PATH so it survives container restarts. The state dir
; is bind-mounted from /volume1/docker/synology-mcp-hub/state/trakt on the
; Synology host.
; PYTHONPATH includes the trakt source root so absolute imports
; (`config.auth`, `models.auth`, etc.) resolve.
; -------------------------------------------------------------------------
[program:trakt]
command=supergateway
    --stdio "python /mcps/trakt/server.py"
    --outputTransport streamableHttp
    --streamableHttpPath /mcp
    --healthEndpoint /healthz
    --port 8123
    --logLevel info
    --cors
directory=/mcps/trakt
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=10
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
environment=PORT="8123",PYTHONPATH="/mcps/trakt",TRAKT_AUTH_TOKEN_PATH="/state/trakt/auth_token.json"

; -------------------------------------------------------------------------
; openfoodfacts-mcp (Open Food Facts public API wrapper)
;
; Native HTTP via the upstream npm package @jagjeevan/openfoodfacts-mcp
; (15-star Node TS, MIT licensed, actively maintained). Installed globally
; in the runtime layer; binary is `off-mcp-server`. Defaults to HTTP
; transport when launched without --stdio. No auth — public APIs only
; (api.openfoodfacts.org + prices.openfoodfacts.org).
; -------------------------------------------------------------------------
[program:openfoodfacts]
command=off-mcp-server
directory=/tmp
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=10
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
environment=PORT="8124",NODE_ENV="production"

; -------------------------------------------------------------------------
; pbs-mcp (Australian Pharmaceutical Benefits Scheme API wrapper)
;
; Native HTTP via Express. Recovered compiled JS from the Cloud Run image
; (sha:e19ba6e2…); 6 ergonomic per-operation tools (pbs_search_medicines,
; pbs_get_pricing, pbs_list_schedules, pbs_get_restrictions,
; pbs_summary_of_changes, pbs_raw_query). Public API — no auth required.
; TRANSPORT=http selects the Express HTTP server vs default stdio.
; -------------------------------------------------------------------------
[program:pbs]
command=node /mcps/pbs/dist/index.js
directory=/mcps/pbs
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=10
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
environment=PORT="8125",NODE_ENV="production",TRANSPORT="http",MCP_ENDPOINT="/mcp"
