# supervisord.conf
; supervisord configuration for the Synology MCP Hub.
;
; Each MCP runs as its own child process on a distinct internal port.
; Logs are streamed to stdout/stderr so `docker logs mcp-hub` shows
; everything (per-program log lines are prefixed with the program name).
;
; Programs:
;   arr           - python stdio MCP wrapped by supergateway (since mcp-arr
;                   hasn't been ported to native HTTP transport)
;   dispatcharr   - python FastMCP, native HTTP transport
;   gramps        - node TS MCP, native HTTP transport
;   trakt         - python stdio MCP wrapped by supergateway
;   openfoodfacts - npm @jagjeevan/openfoodfacts-mcp, native HTTP
;   pbs           - node Express HTTP MCP (recovered compiled JS)
;   whatsupp      - node Express HTTP MCP (medication tracker, Supabase backend)
;   nanobanana    - python FastMCP, native HTTP (Gemini image gen)
;   gws           - @googleworkspace/cli@0.7.0 stdio MCP wrapped by supergateway

[supervisord]
nodaemon=true
loglevel=%(ENV_SUPERVISOR_LOGLEVEL)s
logfile=/dev/stdout
logfile_maxbytes=0
pidfile=/run/supervisord.pid

[unix_http_server]
file=/run/supervisor.sock

[supervisorctl]
serverurl=unix:///run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

; -------------------------------------------------------------------------
; arr-mcp (Sonarr / Radarr / Prowlarr wrapper)
;
; mcp-arr installs a console script `arr-mcp` that runs the stdio MCP.
; supergateway bridges stdio ↔ streamable HTTP and exposes /mcp on 8120.
; Env (SONARR_URL, SONARR_API_KEY, RADARR_*, PROWLARR_*) is supplied by
; docker-compose env_file: env/arr.env.
; -------------------------------------------------------------------------
[program:arr]
command=supergateway
    --stdio "arr-mcp"
    --outputTransport streamableHttp
    --streamableHttpPath /mcp
    --healthEndpoint /healthz
    --port 8120
    --logLevel info
    --cors
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=10
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
stdout_events_enabled=true
stderr_events_enabled=true
environment=PORT="8120"

; -------------------------------------------------------------------------
; dispatcharr-mcp (TV channel manager wrapper)
;
; Native HTTP via FastMCP — no supergateway needed. server.py reads PORT
; and serves /mcp directly. Env (DISPATCHARR_URL, DISPATCHARR_API_KEY,
; STREAM_PROBE_*) supplied by docker-compose.
; -------------------------------------------------------------------------
[program:dispatcharr]
command=python /mcps/dispatcharr/server.py
directory=/mcps/dispatcharr
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=10
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
environment=PORT="8121",MCP_ENDPOINT="/mcp"

; -------------------------------------------------------------------------
; gramps-mcp (Gramps Web genealogy app wrapper)
;
; Native HTTP via Node TS app. dist/index.js reads PORT and serves /mcp.
; Env (GRAMPSWEB_URL, GRAMPS_USERNAME, GRAMPS_PASSWORD) supplied by
; docker-compose.
; -------------------------------------------------------------------------
[program:gramps]
command=node /mcps/gramps/dist/index.js
directory=/mcps/gramps
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=10
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
environment=PORT="8122",MCP_TRANSPORT="http"

; -------------------------------------------------------------------------
; trakt-mcp (Trakt.tv API wrapper)
;
; stdio MCP wrapped by supergateway. Trakt has its own per-user device-code
; OAuth (separate from the Worker's OAuth) — token persists to
; TRAKT_AUTH_TOKEN_PATH so it survives container restarts. The state dir
; is bind-mounted from /volume1/docker/synology-mcp-hub/state/trakt on the
; Synology host.
; PYTHONPATH includes the trakt source root so absolute imports
; (`config.auth`, `models.auth`, etc.) resolve.
; -------------------------------------------------------------------------
[program:trakt]
command=supergateway
    --stdio "python /mcps/trakt/server.py"
    --outputTransport streamableHttp
    --streamableHttpPath /mcp
    --healthEndpoint /healthz
    --port 8123
    --logLevel info
    --cors
directory=/mcps/trakt
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=10
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
environment=PORT="8123",PYTHONPATH="/mcps/trakt",TRAKT_AUTH_TOKEN_PATH="/state/trakt/auth_token.json"

; -------------------------------------------------------------------------
; openfoodfacts-mcp (Open Food Facts public API wrapper)
;
; Native HTTP via the upstream npm package @jagjeevan/openfoodfacts-mcp
; (15-star Node TS, MIT licensed, actively maintained). Installed globally
; in the runtime layer; binary is `off-mcp-server`. Defaults to HTTP
; transport when launched without --stdio. No auth — public APIs only
; (api.openfoodfacts.org + prices.openfoodfacts.org).
; -------------------------------------------------------------------------
[program:openfoodfacts]
command=off-mcp-server
directory=/tmp
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=10
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
environment=PORT="8124",NODE_ENV="production"

; -------------------------------------------------------------------------
; pbs-mcp (Australian Pharmaceutical Benefits Scheme API wrapper)
;
; Native HTTP via Express. Recovered compiled JS from the Cloud Run image
; (sha:e19ba6e2…); 6 ergonomic per-operation tools (pbs_search_medicines,
; pbs_get_pricing, pbs_list_schedules, pbs_get_restrictions,
; pbs_summary_of_changes, pbs_raw_query). Public API — no auth required.
; TRANSPORT=http selects the Express HTTP server vs default stdio.
; -------------------------------------------------------------------------
[program:pbs]
command=node /mcps/pbs/dist/index.js
directory=/mcps/pbs
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=10
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
environment=PORT="8125",NODE_ENV="production",TRANSPORT="http",MCP_ENDPOINT="/mcp"

; -------------------------------------------------------------------------
; whatsupp-mcp (medication & supplement tracker)
;
; Native HTTP via Express, using StreamableHTTPServerTransport from the
; @modelcontextprotocol/sdk. dist/index.js reads PORT + MCP_ENDPOINT and
; serves /mcp directly. /health is also exposed for the outer healthcheck.
;
; Backed by Supabase Postgres (project rafybjboggundidsdwjz). Auth: the
; service_role key (RLS-bypass) — same pattern as the NFC pathway from
; HA Balmoral, since RLS policies on the meds tables require a row in
; public.patient_access for auth.uid() and patient_access is empty.
;
; Migrated 2026-05-10 from Cloud Run (australia-southeast1/whatsupp-mcp-server)
; to consolidate the fleet on the Hub + drop the URL-token auth pattern.
;
; Env (supplied via env_file: env/whatsupp.env):
;   - SUPABASE_URL
;   - SUPABASE_SERVICE_ROLE_KEY
; -------------------------------------------------------------------------
[program:whatsupp]
command=node /mcps/whatsupp/dist/index.js
directory=/mcps/whatsupp
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=10
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
environment=PORT="8129",NODE_ENV="production",TRANSPORT="http",MCP_ENDPOINT="/mcp"

; -------------------------------------------------------------------------
; nanobanana-mcp (Gemini image generation wrapper)
;
; PyPI package zhongweili/nanobanana-mcp-server, installed as `nanobanana-mcp-server`
; console script. Native HTTP via FastMCP — reads FASTMCP_TRANSPORT/HOST/PORT
; from env to bind on :8127. Generates images via Gemini Flash 3.1 (default),
; Pro 3, or 2.5 Flash (legacy). Returns inline base64 PNG content blocks so
; Claude.ai web renders previews directly.
;
; Env contract:
;   - FASTMCP_TRANSPORT=http              (selects HTTP server; default would be stdio)
;   - FASTMCP_HOST=0.0.0.0                (bind on all interfaces — default 127.0.0.1)
;   - FASTMCP_PORT=8127                   (matches port mapping)
;   - IMAGE_OUTPUT_DIR=/tmp/nanobanana-images  (ephemeral; images are also returned inline)
;   - GEMINI_API_KEY=...                  (supplied via env_file: env/nanobanana.env)
;
; Image processing peaks ~700-900 MiB during multi-image renders. Synology
; has 18 GB physical RAM so this is comfortably absorbed (~5%).
; -------------------------------------------------------------------------
[program:nanobanana]
command=nanobanana-mcp-server
directory=/tmp
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=10
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
environment=FASTMCP_TRANSPORT="http",FASTMCP_HOST="0.0.0.0",FASTMCP_PORT="8127",IMAGE_OUTPUT_DIR="/tmp/nanobanana-images"

; -------------------------------------------------------------------------
; gws-mcp (Google Workspace CLI MCP — Drive / Tasks / Gmail / Contacts)
;
; Wraps the official @googleworkspace/cli@0.7.0 npm package. The CLI ships an
; `mcp` subcommand that exposes ~120 tools across the requested services
; (drive ~28, tasks ~14, gmail ~50, people ~24). The CLI itself only speaks
; stdio, so supergateway bridges to streamable-HTTP on port 8128.
;
; mcps/gws/start.sh hydrates ~/.config/gws/ from the GWS_CREDENTIALS env var
; (unencrypted refresh-token JSON, output of `gws auth export --unmasked`)
; before exec'ing `gws mcp -s drive,tasks,gmail,people`. Three Cloud Run
; lessons preserved here: ca-certificates installed in runtime layer (gws Go
; binary needs system CAs); ~/.config/gws/ needs FOUR files synthesized, not
; just credentials.json; supergateway's child_process.spawn doesn't reliably
; inherit env so start.sh re-exports HOME + KEYRING_BACKEND.
;
; gws CLI version is PINNED to 0.7.0 — `gws mcp` subcommand was removed in
; 0.8.0 (upstream PR #275). Bumping breaks the MCP. See Dockerfile comment.
;
; Env (supplied via env/gws.env on Synology host):
;   - GWS_CREDENTIALS=<JSON>  (unencrypted refresh token + client_id/secret)
; -------------------------------------------------------------------------
[program:gws]
command=supergateway
    --stdio "/mcps/gws/start.sh"
    --outputTransport streamableHttp
    --streamableHttpPath /mcp
    --healthEndpoint /healthz
    --port 8128
    --logLevel info
    --cors
directory=/mcps/gws
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=10
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
environment=PORT="8128"
